# sein/templates/default/seinworkers.conf.erb
[program:replicator]
process_name=%(program_name)s_%(process_num)02d
command="/usr/local/zend/bin/php /opt/sein/current/vendor/bin/zf.php start worker sein_replicator"
numprocs=2
autorestart=true
startretries=2
environment=ZF_CONFIG_FILE="/opt/sein/current/application/configs/zf.ini",ZEND_TOOL_INCLUDE_PATH="/opt/sein/current/vendor/zendframework/zendframework1/library:/opt/sein/current/providers"

[program:forecast]
process_name=%(program_name)s_%(process_num)02d
command="/opt/sein/current/vendor/bin/zf.php start worker sein_forecast"
numprocs=6
autorestart=true
startretries=2
environment=ZF_CONFIG_FILE="/opt/sein/current/application/configs/zf.ini",ZEND_TOOL_INCLUDE_PATH="/opt/sein/current/vendor/zendframework/zendframework1/library:/opt/sein/current/providers"

[group:seinworkers]
programs=replicator,forecast